# luci-theme-opentomcat/root/etc/uci-defaults/30_luci-theme-opentomcat
#!/bin/sh
if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
	set luci.themes.Opentomcat=/luci-static/opentomcat
	commit luci
	EOF
fi

exit 0
