# luci-theme-argon-dark/root/etc/uci-defaults/30_luci-theme-argon-dark
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Argon_Dark=/luci-static/argon-dark
	commit luci
EOF
exit 0
